# src/extends/res/websys/tc/example.cron
30 01 * * *     (ls /usr/local  >> /tmp/crontab_tc.log 2>&1)
